# scripts/com.teamlift.banter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.teamlift.banter</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/simongilmore/Documents/Projects/team-lift/scripts/refresh-banter.sh</string>
  </array>
  <!-- Six-hourly at :47 past; launchd runs a missed slot on wake from sleep.
       The script itself hashes each dashboard section and only calls out to
       claude when something actually changed since the last run. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>0</integer>
      <key>Minute</key><integer>47</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>6</integer>
      <key>Minute</key><integer>47</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>12</integer>
      <key>Minute</key><integer>47</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>47</integer>
    </dict>
  </array>
  <!-- Also fire at login, covering slots the Mac was powered off for.
       The script's own hash-based change detection keeps this cheap. -->
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
