# scripts/com.teamlift.banter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.teamlift.banter</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/simongilmore/Projects/team-lift/scripts/refresh-banter.sh</string>
  </array>
  <!-- Hourly; launchd coalesces missed runs into one on wake from sleep.
       Cheap to run this often: the script hashes each dashboard section and
       only calls out to claude when something actually changed, so a quiet
       hour costs one date-bump PATCH and no AI call. -->
  <key>StartInterval</key>
  <integer>3600</integer>
  <!-- Also fire at login, covering slots the Mac was powered off for.
       The script's own hash-based change detection keeps this cheap. -->
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
