# scripts/com.teamlift.banter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.teamlift.banter</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/simongilmore/Documents/Projects/team-lift/scripts/refresh-banter.sh</string>
  </array>
  <!-- Daily at 6:47am; launchd runs a missed slot on wake from sleep. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>47</integer>
  </dict>
  <!-- Also fire at login, covering days the Mac was powered off at 6:47.
       The script exits early if today's banter already exists. -->
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
